# etc/supervisor.conf
[program:testapp]
command = /home/eli/testapp/venv/bin/gunicorn app:app -b localhost:8000
directory = /home/eli/testapp
user = eli
