# etc/supervisor.conf
[program:testapp]
command = /home/eli/testapp/venv/bin/gunicorn app:app --bind unix:/tmp/testapp.sock --log-file -
directory = /home/eli/testapp
user = eli
